# pim/.config/systemd/user/khal-reminder.service
[Unit]
Description=Khal reminder daemon

[Service]
Type=oneshot
ExecStart=%h/auto-scripts/khal-reminder-daemon.sh

